# etc/cron.d/sisiya-client-checks
### SisIYA : Client checks.
*/10 * * * * root /usr/bin/perl -I/etc/sisiya/sisiya-client-checks /usr/share/sisiya-client-checks/utils/sisiya_check.pl 15
