# etc/cron.d/sisiya-check-expired
#### SisIYA : Check for expired messages.
0-59/5 * * * * root php -q __WEB_BASE_DIR__/utils/sisiya_check_expired.php __WEB_BASE_DIR__
