# etc/cron.d/sisiya-remote-checks
#### SisIYA : Remote checks.
0-59/5 * * * * root /usr/bin/perl -I/etc/sisiya/sisiya-remote-checks -I/etc/sisiya/sisiya-client-checks /usr/share/sisiya-remote-checks/utils/sisiya_remote_check.pl 10 
