# etc/cron.d/sisiya-rss
#### SisIYA : Generate the RSS xml file
* * * * * root php -q __WEB_BASE_DIR__/utils/sisiya_rss.php > __WEB_BASE_DIR__/utils/sisiya_rss.xml
