# etc/cron.d/sisiya-rss
#### SisIYA : Generate the RSS xml file
* * * * * root php -q /usr/share/sisiya-webui-php/utils/sisiya_rss.php > /usr/share/sisiya-webui-php/utils/sisiya_rss.xml
