# etc/cron.d/sisiya-check-expired
#### SisIYA : Check for expired messages.
* * * * * root php -q /usr/share/sisiya-webui-php/utils/sisiya_check_expired.php /usr/share/sisiya-webui-php /usr/share/sisiya-client-checks/utils/sisiya_send_message2.sh
