# etc/cron.d/sisiya_alerts
#### SisIYA : Send alerts if any.
* * * * * root php /var/www/html/sisiya/sisiya_alerts.php > /dev/null
