# etc/cron.d/sisiya-archive
#### SisIYA : Update the SisIYA database and archive the system history messages
02 2 * * * root php -q __WEB_BASE_DIR__/utils/sisiya/sisiya_updatedb.php __WEB_BASE_DIR__
