# etc/cron.d/sisiya_check_expired
#### SisIYA : Check expired messages.
0-59/5 * * * * root /opt/sisiya-server-checks/scripts/sisiya_check_expired_mysql.sh /opt/sisiya-server-checks/conf/sisiya_server_checks.conf
