# etc/cron.d/sisiya_client_checks_perl
### SisIYA : Cleint checks.
0-59/5 * * * * root perl -I/opt/sisiya-client-checks /opt/sisiya-client-checks/utils/sisiya_check.pl 10 2>/dev/null
