# etc/cron.d/sisiya-alerts
#### SisIYA : Send alerts if any.
* * * * * root php -q /usr/share/sisiya-webui-php/utils/sisiya_alerts.php /usr/share/sisiya-webui-php > /dev/null
